# x-viredero@.service
[Service]
ExecStart=/usr/bin/start_x-viredero.sh %I
